# cert/renew-cert.service
[Unit]
Description=Renew self-signed TLS certificate

[Service]
Type=oneshot
ExecStart=/usr/local/bin/renew-cert.sh
StandardOutput=journal
StandardError=journal